# hw_button_scan.service.in
[Unit]
Description=key_control

[Service]
Type=simple
Restart=always
RestartSec=5
ExecStart=${LOC}/venv/bin/python3 ${LOC}/hw_button_scan.py
ExecStop=/bin/systemctl kill hw_button_scan

[Install]
WantedBy=multi-user.target


 
